# workers/crontab
*/1 * * * * /Users/hackreactor/code/krgeppert/2013-06-web-historian/workers/htmlfetcher.js -update -config=myconfig